# xray-init
#!/bin/sh /etc/rc.common
START=99
STOP=10
USE_PROCD=1
PROG=/usr/local/bin/xray
CONFIG=/usr/local/etc/xray/config.json

start_service() {
    mkdir -p /var/log/xray
    /usr/local/bin/xray-tproxy-setup.sh
    procd_open_instance
    procd_set_param command $PROG run -c $CONFIG
    procd_set_param env XRAY_LOCATION_ASSET=/usr/local/etc/xray
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile /var/run/xray.pid
    procd_set_param limits nofile="65536 65536"
    procd_close_instance

    # OOM protection: после старта снижаем oom_score_adj у xray до -500 —
    # OOM-killer будет убивать почти все остальные процессы (logd, dnsmasq,
    # уведомительные скрипты) раньше xray. xray — самое важное.
    (
        for i in 1 2 3 4 5 6 7 8 9 10; do
            sleep 1
            PID=$(cat /var/run/xray.pid 2>/dev/null)
            if [ -n "$PID" ] && [ -f "/proc/$PID/oom_score_adj" ]; then
                echo -500 > /proc/$PID/oom_score_adj 2>/dev/null && break
            fi
        done
    ) &
}

stop_service() {
    nft delete table inet xray_tproxy 2>/dev/null
    ip rule del fwmark 1 lookup 100 2>/dev/null
}

reload_service() {
    stop_service
    sleep 1
    start_service
}
